# supervisord.conf
[supervisord]
nodaemon=true

[program:python]
command=python3 /Redstoneguy129/update.py

[program:nginx]
command=nginx -c /etc/nginx/nginx.conf
